# infrastructure/hermes/ansible/roles/hermes/templates/hermes.service.j2
[Unit]
Description=Hermes AI Agent
After=network.target

[Service]
Type=simple
User={{ hermes_user }}
Group={{ hermes_group }}
WorkingDirectory={{ hermes_app_dir }}
ExecStart={{ hermes_venv_dir }}/bin/gunicorn "lib.interfaces.api:create_app()" \
    --bind 0.0.0.0:{{ hermes_port }} --workers 1 --timeout 300
Environment=HERMES_API_TOKEN={{ hermes_api_token }}
Environment=PYTHONPATH={{ hermes_app_dir }}
Restart=on-failure
RestartSec=3

# Resource limits
LimitNOFILE={{ max_open_files }}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ hermes_log_dir }} {{ hermes_config_dir }}

[Install]
WantedBy=multi-user.target
